# packaging/paneld.service
# ttlcd-panel daemon — systemd USER service (canonical unit).
#
# The `ttlcd` CLI installs this automatically: `ttlcd up` copies it to
# ~/.config/systemd/user/, enables, and starts it. Manage with
# `ttlcd up|down|restart|status|logs`. This file is the single source of truth
# for the unit definition.
#
# Manual install (equivalent to `ttlcd up`):
#   mkdir -p ~/.config/systemd/user
#   cp /home/joe/projects/ttlcd/packaging/paneld.service ~/.config/systemd/user/
#   systemctl --user daemon-reload && systemctl --user enable --now paneld
#
# Run without an active login session (headless / after reboot): loginctl enable-linger joe

[Unit]
Description=Thermaltake LCD panel daemon (ttlcd)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/joe/projects/ttlcd/.venv/bin/paneld --idle system --rotate 5
WorkingDirectory=/home/joe/projects/ttlcd
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
